# main.cron
* 0 1-31 * * /usr/bin/python s_20/main.py
